# docker/supervisor.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true

[program:queue]
command=php /var/www/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
